# src/modules/octopi/filesystem/root/etc/systemd/system/webcamd.service
[Unit]
Description=the OctoPi webcam daemon with the user specified config
ConditionPathExists=/etc/octopi_streamer/mjpeg

[Service]
WorkingDirectory=/root/bin
ExecStart=/root/bin/webcamd
Restart=always
Type=forking

[Install]
WantedBy=multi-user.target
